# crontab
*/1 * * * * cd /py_cronjob && /usr/bin/python3.10 main.py >> cron_log.log 2>&1
