# crontab
#0 17-23/3,0-7/3 * * * cd /py_cronjob && python main.py >> cron_log.log 2>&1
0 9-22/3 * * * cd /py_cronjob && python main.py >> cron_log.log 2>&1
